# ocr-server/ocr-server.service
[Unit]
Description=Surya OCR Server (GPU)
After=network.target

[Service]
Type=simple
WorkingDirectory=PLACEHOLDER_HOME/local-ai-search/ocr-server
ExecStart=PLACEHOLDER_HOME/local-ai-search/ocr-server/venv/bin/uvicorn server:app --host 0.0.0.0 --port 8090

# AMD GPU stability
Environment=HSA_ENABLE_SDMA=0
Environment=AMD_SERIALIZE_KERNEL=3
Environment=ROCM_FORCE_DISABLE_LAZY_ALLOC=1
Environment=TORCH_ROCM_AOTRITON_ENABLE_EXPERIMENTAL=1
Environment=RECOGNITION_BATCH_SIZE=64
Environment=DETECTOR_BATCH_SIZE=8
Environment=TORCH_DEVICE=cuda
Environment=PATH=PLACEHOLDER_HOME/local-ai-search/ocr-server/venv/bin:/opt/rocm/bin:/usr/local/bin:/usr/bin:/bin

# Auto-restart on crash
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
